# deploy/enclume-fusion-client.service
[Unit]
Description=Enclume Shared Integration Vite Client
After=network.target enclume-fusion-server.service

[Service]
Type=simple
User=codex
WorkingDirectory=/home/codex/Enclume-fusion/client
EnvironmentFile=/home/codex/Enclume-fusion/.env
ExecStart=/usr/bin/node /home/codex/Enclume-fusion/client/node_modules/vite/bin/vite.js --host 0.0.0.0 --port 8393 --strictPort
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
